# vps/hermes-api.service
[Unit]
Description=Hermes API shim (website -> Hermes Agent CLI)
After=network.target

[Service]
Type=simple
User=root
Environment=HOME=/root
EnvironmentFile=/etc/hermes-api.env
ExecStart=/usr/bin/python3 /opt/hermes-api/hermes_api.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
